# deploy/production/pulsebrief-backup.timer
[Unit]
Description=Run PulseBrief backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
